# iso/skeleton/etc/profile.d/smgl-install.sh
#!/bin/bash

# This includes the SMGL installer functions
. /usr/share/smgl.install/bashrc
